# homelab/openclaw/openclaw.container
[Unit]
Description=OpenClaw gateway
Wants=network-online.target
Wants=homelab-network.service
After=network-online.target
After=homelab-network.service

[Container]
ContainerName=openclaw
Image=docker.io/alpine/openclaw:latest
EnvironmentFile=%h/.config/containers/systemd/openclaw.env
UserNS=keep-id
Volume=__OPENCLAW_DATA_DIR__:/home/node/.openclaw:Z
Network=homelab.network
PublishPort=18789:18789
PublishPort=18790:18790
AutoUpdate=registry

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target
